# debian/postrm
#!/bin/sh
set -e

if [ purge = "$1" ]; then
   update-rc.d dhcp-helper remove >/dev/null
fi
